# systemd/fetch-next-event.timer
[Unit]
Description=Fetch next event regularly

[Timer]
#OnCalendar=*:0/15
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
