# systemd/fetch-next-event.service
[Unit]
Description=Fetch next event

[Service]
ExecStart=/home/manfred/.bin/fetch-next-event.sh
User=manfred

[Install]
WantedBy=multi-user.target
