# systemd/fetch-next-event.service
[Unit]
Description=Fetch next event using gcalcli

[Service]
Type=oneshot
ExecStart=/home/manfred/.bin/fetch-next-event
