# home/.config/service/pipewire/run
#!/bin/zsh

source $HOME/.config/zsh/.zshrc
exec pipewire
