# home/.config/service/xbindkeys/run
#!/bin/sh
exec xbindkeys --file .xbindkeysrc --nodaemon
