# home/.config/service/pipewire-pulse/run
#!/bin/zsh

SVDIR=$HOME/.config/service sv check pipewire || exit 1

source $HOME/.config/zsh/.zshrc
exec pipewire-pulse
