# home/.config/service/pulseaudio/run
#!/bin/sh

exec pulseaudio --system=false --daemonize=no --exit-idle-time=-1
