# home/.config/service/xss-lock/run
#!/bin/sh
xset s on
xset s 60
exec xss-lock slock +resetsaver
