# etc/sv/runsvdir-manfred/run
#!/bin/sh

export USER="manfred"
export HOME="/home/manfred"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.config/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
